# packer/crontab
*/5 * * * * /usr/local/bin/reload_ssh_keys >/dev/null 2>&1
